# layout/DEBIAN/postinst
#!/bin/sh

launchctl load /Library/LaunchDaemons/com.bishopfox.sshd.plist
launchctl start com.bishopfox.sshd

